# ops/systemd/mme-engine.service
# market-microstructure-engine — 24/5 data-collection service (read-only, NO trading).
#
# Production path for long-running data collection. Use this INSTEAD of the cron
# lifecycle (ops/crontab.example) — never both, they would fight over the process.
#
# Install everything (engine + health + cleanup timers) with:
#   sudo ops/install_systemd.sh
# Or manually:
#   sudo cp ops/systemd/*.service ops/systemd/*.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now mme-engine mme-health.timer mme-cleanup.timer
#   systemctl status mme-engine ; journalctl -u mme-engine -f
#
# MT5-only: the engine pulls quotes from the read-only MT5 bridge sidecar
# (default 127.0.0.1:7777). Make sure the bridge is running (under Wine) — the
# engine pings it at startup, refuses any non-demo account, and retries per cycle.
# Override the bridge endpoint with Environment=MME_MT5_HOST=.. MME_MT5_PORT=..

[Unit]
Description=Market Microstructure Engine (24/5 data collection, read-only, no trading)
After=network-online.target
Wants=network-online.target
# Crash-loop protection MUST live in [Unit] (not [Service]): if the engine fails
# 5 times within 300s, stop retrying and enter 'failed' -> OnFailure pages us.
StartLimitIntervalSec=300
StartLimitBurst=5
OnFailure=mme-alert.service

[Service]
Type=simple
User=kadir70
WorkingDirectory=/home/kadir70/market-microstructure-engine/build
EnvironmentFile=-/home/kadir70/market-microstructure-engine/ops/.env
ExecStart=/home/kadir70/market-microstructure-engine/build/engine

# Auto-restart on any non-clean exit, with backoff so a crash-loop can't hammer
# the API (and burn the free-tier daily quota).
Restart=always
RestartSec=10

# Clean shutdown: engine handles SIGTERM (interruptibleSleep -> ~1s exit).
KillSignal=SIGTERM
TimeoutStopSec=20

# Bound journald growth from the engine's per-cycle stdout over weeks of running.
# (The authoritative log is still the rotating ../logs/engine.log, 10MB x 7.)
LogRateLimitIntervalSec=30
LogRateLimitBurst=20

# Light hardening — only needs to read config and write data/ + logs/.
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
